# cdr-exporter.cron
25,55 * * * *   root    if /usr/sbin/ngcp-check_active -q; then /usr/sbin/cdr-exporter >/dev/null; fi
