# cdr-exporter.cron
25 * * * *   root    /usr/sbin/cdr-exporter >/dev/null
55 * * * *   root    /usr/sbin/cdr-exporter >/dev/null
