# event-exporter.cron
0 0 * * *   root    /usr/sbin/event-exporter >/dev/null
