# event-exporter.cron
0 0 * * *   root    if /usr/sbin/ngcp-check_active -q; then /usr/sbin/event-exporter >/dev/null; fi
